# ami/resources/etc/profile.d/sh.local
# Ensure /usr/local/bin is at the start of the path
# ex: sw=4 ts=4 et ai

case ":${PATH}:" in
    :/usr/local/bin:*)  ;;
    *:/usr/local/bin:*) PATH="/usr/local/bin:"`echo "$PATH" | sed -e 's!:/usr/local/bin!!'`;;
    *)                  PATH="/usr/local/bin:$PATH";;
esac
